# cookbooks/scintilla/templates/default/daphne.conf.erb
[program: daphne]
command=sudo daphne <%= @app_folder %>.asgi:channel_layer --port <%= @port %> --bind localhost
directory=/var/django/<%= @app_folder %>/
user=<%= @user %>
autostart=True
autorestart=True
redirect_stderr=True
environment=<%= @supervisor_environment %>