# cookbooks/scintilla/templates/default/worker.conf.erb
[program: worker]
command=python3 manage.py runworker
directory=/var/django/<%= @app_folder %>/
user=<%= @app_user %>
autostart=True
autorestart=True
redirect_stderr=True
environment=<%= @supervisor_environment %>
